# etc/rc.local

#-- turn on NumLock for all terminals
INITTY=/dev/tty[1-8]
for tty in $INITTY; do
    setleds -D _num < $tty
done


exit 0;
